# scripts/docker-firewall.service
[Unit]
Description=Docker firewall rules (block external access to 8000/8001)
After=docker.service network-online.target
Requires=docker.service
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -c '\
    EXT_IFACE=$(ip route get 8.8.8.8 | grep -oP "dev \\K\\S+" | head -1) && \
    iptables -D DOCKER-USER -i $EXT_IFACE -p tcp --dport 8000 -j DROP 2>/dev/null; \
    iptables -D DOCKER-USER -i $EXT_IFACE -p tcp --dport 8001 -j DROP 2>/dev/null; \
    iptables -I DOCKER-USER -i $EXT_IFACE -p tcp --dport 8000 -j DROP && \
    iptables -I DOCKER-USER -i $EXT_IFACE -p tcp --dport 8001 -j DROP'

[Install]
WantedBy=multi-user.target
